# config/autostart/copyq.desktop
[Desktop Entry]
Exec=/usr/bin/copyq
Name=copyq
Type=Application
