# config/autostart/numlockx.desktop
[Desktop Entry]
Exec=/usr/bin/numlockx on
Name=numlockx
Type=Application
